# zenvx-iso/config/includes.chroot/etc/systemd/system/zenvx-log-manager.timer
[Unit]
Description=ZenvX Log Manager Timer

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
Persistent=yes

[Install]
WantedBy=timers.target
